# Dots-WM/Bspwm-Nitro-Artix/etc/rc.local
#!/bin/bash -e
###  Chmod 755


$HOME/.bin/Dual-DM.sh

exit 0

